# launchd/com.tracker-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tracker-web</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/tracker</string>
        <string>serve</string>
        <string>--no-browser</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/tracker-web.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tracker-web.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>TRACKER_VIEWS</key>
        <string>dashboard,daily,backlog,topics,meetings</string>
    </dict>
</dict>
</plist>
